# supervisor/tomcat_svd.conf
[program:tomcat]
command=$CATALINA_HOME/bin/catalina.sh run