# roles/ContainerizeGeo/files/tomcat.conf
[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/tomcat.log
redirect_stderr=true