# dockerfile/context/tomcat_svd.conf
[program:tomcat]
command=catalina.sh run